# service.erb
#!/bin/sh
# Init script for <%= @name %> 
# This file is generated by Chef
# Changes will be overwritten
#
### BEGIN INIT INFO
# Provides: <%= @name %>
# Default-Start: 2 3 4 5
# Default-Stop: 0 1 6
# Description: Init script for <%= @name %>
### END INIT INFO

. /lib/lsb/init-functions

_pid(){
  local p=`pgrep -f <%= @process %>`
  echo "$p"
}

_start(){
    pid=`_pid`
    if [ $pid ] ; then
      echo "<%= @name %> is already running. pid: $pid"
    else
      echo -n "Starting <%= @name %>:"
      su -s /bin/bash -c "nohup <%= @init_script %> <%= @config %> > <%= @log_dir %>/<%= @name %>.out 2> <%= @log_dir %>/<%= @name %>.err &" <%= @user %>
    fi
}

_stop(){
    pid=`_pid`
    if ! [ $pid ] ; then
      echo "<%= @name %> is not running"
    else
      echo -n "Stopping <%= @name %>:"
      su -s /bin/bash -c "ps -ef | grep <%= @process %> | grep -v grep | awk '{print \$2}' | xargs kill -SIGTERM" <%= @user%>
    fi
}

_status(){
    pid=`_pid`
    if [ "$pid" = "" ] ; then
        echo "Stopped"
        exit 2
    else
        echo "Running $pid"
        exit 0
    fi
}

log_end_msg() {
  if [ "$1" = 0 ]; then
    echo " [  OK  ]"
  else
    echo " [FAILED]"
  fi
}

case "$1" in
  start)
    _start
    log_end_msg "$?"
  ;;

  stop)
    _stop
    log_end_msg "$?"
  ;;


  restart)
    _stop
    sleep 3s
    _start
    log_end_msg "$?"
  ;;

  status)
  _status
  ;;

  *)
    echo "Usage: /etc/init.d/<%= @name %> {start|stop|restart|status}"
    exit 1
esac